# skyratshotspot-connect.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/connect.sh

[Install]
WantedBy=default.target
